# deploy/systemd/clonr-serverd.service
[Unit]
Description=clonr Server
After=network.target dev-pts.mount
Wants=network.target dev-pts.mount

[Service]
Type=simple

# ---------------------------------------------------------------------------
# Startup
# ---------------------------------------------------------------------------
# Ensure devpts is mounted before starting — required for systemd-nspawn shell sessions.
ExecStartPre=/bin/bash -c "mountpoint -q /dev/pts || mount -t devpts devpts /dev/pts -o ptmxmode=666,mode=620,gid=5"
ExecStart=/usr/local/bin/clonr-serverd --pxe

# ---------------------------------------------------------------------------
# Environment
# ---------------------------------------------------------------------------
Environment=CLONR_LISTEN_ADDR=0.0.0.0:8080
Environment=CLONR_IMAGE_DIR=/var/lib/clonr/images
Environment=CLONR_DB_PATH=/var/lib/clonr/db/clonr.db
Environment=CLONR_BOOT_DIR=/var/lib/clonr/boot
Environment=CLONR_TFTP_DIR=/var/lib/clonr/tftpboot
Environment=CLONR_PXE_ENABLED=true
Environment=CLONR_PXE_INTERFACE=eth1
Environment=CLONR_PXE_RANGE=10.99.0.100-10.99.0.200
Environment=CLONR_PXE_SERVER_IP=10.99.0.1
Environment=TMPDIR=/var/lib/clonr/tmp

# ---------------------------------------------------------------------------
# Restart policy
# ---------------------------------------------------------------------------
Restart=on-failure
RestartSec=5s

# ---------------------------------------------------------------------------
# Resource limits
# ---------------------------------------------------------------------------
# 8G ceiling: clonr-serverd is a combined management plane + ISO builder that
# spawns QEMU children. Each QEMU builder VM requests -m 2048 (2G), and
# multiple concurrent builds plus Go runtime overhead easily exceed 2G.
# The old 2G limit was appropriate for a pure control plane; with builders
# co-located it caused SIGKILL at mid-build (confirmed: Rocky 10 SeaBIOS
# OOM'd at ~47s during package download).
#
# Follow-on: introduce clonr-builders.slice so builder QEMU children are
# accounted separately from the management plane cgroup. The code change to
# invoke `systemd-run --scope --slice=clonr-builders.slice` will land in a
# subsequent Dinesh commit. For now, 8G on the parent unit unblocks builds.
MemoryMax=8G

# Deprioritise OOM kill — clonr-serverd is the control plane; let other
# processes be killed first under memory pressure. Range: -1000..1000.
OOMScoreAdjust=-500

# ---------------------------------------------------------------------------
# Security hardening
# ---------------------------------------------------------------------------
# NoNewPrivileges is intentionally set to false here. Setting it to true would
# be inherited by ALL child processes across fork()/exec() — including the
# systemd-nspawn processes launched via "systemd-run --scope" for shell
# sessions. The NoNewPrivileges flag is a per-process Linux security bit set
# via prctl(PR_SET_NO_NEW_PRIVS, 1); it is inherited and cannot be unset by
# child processes. Even when nspawn runs in a separate cgroup scope, it still
# inherits the flag from the clonr-serverd parent, causing pivot_root(2) to
# fail with EPERM ("Failed to move root directory: Operation not permitted").
#
# Security is still enforced by CapabilityBoundingSet and AmbientCapabilities
# below, which tightly restrict which capabilities are available. The process
# and its children cannot gain capabilities beyond what is listed there.
NoNewPrivileges=false

# ProtectSystem: set to "full" rather than "strict" for the ISO build pipeline.
# "strict" creates a private mount namespace where the root filesystem is
# read-only and the util-linux mount(2) wrapper enforces /etc/fstab for all
# block device mounts — even with CAP_SYS_ADMIN — because libmount detects
# the restricted namespace and treats the process as non-privileged.
# This blocks losetup-based rootfs extraction after QEMU install.
#
# "full" makes /usr and /boot read-only but leaves / writable, which keeps
# the system-level protections without the private-namespace restrictions
# that break the mount utility. The follow-on fix is for Dinesh to replace
# exec.Command("mount", ...) with unix.Mount() in the extract pipeline, after
# which we can restore ProtectSystem=strict.
ProtectSystem=full

# Prevent access to /home, /root.
ProtectHome=true

# Private /tmp and /var/tmp — isolates ISO build artifacts from other services.
# NOTE: TMPDIR is set to /var/lib/clonr/tmp above, NOT /tmp, so clonr-serverd
# writes ISO build files to the data volume, not the private /tmp. PrivateTmp
# is safe here.
PrivateTmp=true

# Block access to kernel tunables (/proc/sys, /sys).
ProtectKernelTunables=true

# Do not prevent loading kernel modules. ProtectKernelModules=true implicitly
# drops CAP_SYS_MODULE and also prevents loop partition mounts from working
# correctly on some kernel versions (the seccomp filter it applies blocks
# finit_module). The loop module is pre-loaded via modules-load.d so runtime
# loading is not needed, but the directive itself still interferes with mounts.
# Keeping this false; loop module management is handled by
# /etc/modules-load.d/clonr-loop.conf and /etc/modprobe.d/clonr-loop.conf.
ProtectKernelModules=false

# nspawn shell sessions (spawned via systemd-run --scope --slice=clonr-shells.slice)
# must create a sub-cgroup for the container process. The service's mount namespace
# makes /sys/fs/cgroup read-only when ProtectControlGroups=true, which causes
# nspawn to fail with "Failed to create ... subcgroup: Read-only file system".
# Setting false allows child processes in the same mount namespace to manage their
# own cgroups. The shells slice itself limits resource usage.
ProtectControlGroups=false

# Explicit read-write paths required by clonr-serverd:
#   /var/lib/clonr — all runtime data (DB, images, tftpboot, boot, tmp, iso-cache)
#   /var/log/clonr — log files (may be written by future structured log config)
#   /tmp           — systemd-nspawn needs a real /tmp in some paths; also Go runtime
#   /run           — dnsmasq lease files, PID files if any
#   /dev/pts       — devpts for nspawn pty allocation (ExecStartPre mounts this)
ReadWritePaths=/var/lib/clonr /var/log/clonr /tmp /run /dev/pts

# ---------------------------------------------------------------------------
# Capabilities
# ---------------------------------------------------------------------------
# CAP_NET_BIND_SERVICE  — bind privileged ports (DHCP/67, TFTP/69)
# CAP_SYS_ADMIN         — required by systemd-nspawn for mount namespace creation
#                         and by qemu-kvm for /dev/kvm access; also required for
#                         mount(2) syscall to work without fstab constraints
# CAP_MKNOD             — required by systemd-nspawn to create device nodes inside
#                         the container namespace during shell sessions
# CAP_DAC_OVERRIDE      — util-linux libmount checks whether the bounding set is
#                         restricted to determine if the caller is "fully privileged".
#                         Without a sufficiently broad bounding set, mount enforces
#                         fstab restrictions even for UID 0 + CAP_SYS_ADMIN. Adding
#                         DAC_OVERRIDE and FOWNER to the bounding set satisfies the
#                         libmount privilege check and allows mounting loop partitions
#                         during rootfs extraction without a code change to use the
#                         raw mount(2) syscall instead.
# CAP_DAC_READ_SEARCH   — needed for overlay path traversal during rootfs copy;
#                         allows reading directories the process doesn't own and
#                         bypasses DAC read permission checks on directory traversal.
# CAP_FOWNER            — see CAP_DAC_OVERRIDE note above
# CAP_CHOWN             — required by rsync to preserve ownership during rootfs copy
# CAP_SETUID / CAP_SETGID — required by rsync to restore uid/gid of files
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN CAP_MKNOD CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN CAP_MKNOD CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_SETUID CAP_SETGID

# ---------------------------------------------------------------------------
# Namespace policy
# ---------------------------------------------------------------------------
# systemd-nspawn (used for image shell sessions) requires the ability to create
# mount and user namespaces. RestrictNamespaces=true would block it and break
# the shell feature. We leave namespace creation unrestricted.
RestrictNamespaces=false

# MountFlags=shared: ProtectSystem=strict creates a private mount namespace
# with MS_PRIVATE propagation. This blocks clonr-serverd from mounting loop
# devices for rootfs extraction (mount(2) returns EPERM even with CAP_SYS_ADMIN
# because the namespace itself is locked down). Setting shared propagation allows
# the service to perform mounts within its namespace for losetup/extract workflows.
MountFlags=shared

# ---------------------------------------------------------------------------
# Device access
# ---------------------------------------------------------------------------
# ISO build pipeline uses qemu-kvm which requires /dev/kvm.
# Rootfs extraction after install uses losetup which requires /dev/loop* block
# devices AND /dev/loop-control. Important distinction:
#   block-loop       — covers major 7 block devices (/dev/loop0, /dev/loop1p1,
#                      etc.). Allows read/write/mknod on all loop block nodes.
#   /dev/loop-control — character device (major 10, minor 237) used by
#                       losetup -f to find and allocate a free loop device.
#                       NOT covered by "block-loop". Without this entry,
#                       losetup fails with "cannot find an unused loop device"
#                       because the cgroup device controller denies open(2) on
#                       the control node even for root + CAP_SYS_ADMIN.
# Both entries are required for the full losetup → mount → extract workflow.
DevicePolicy=auto
DeviceAllow=/dev/kvm rwm
DeviceAllow=/dev/net/tun rwm
DeviceAllow=block-loop rwm
DeviceAllow=/dev/loop-control rw

[Install]
WantedBy=multi-user.target
