# deploy/systemd/clonr-autodeploy.service
[Unit]
Description=clonr auto-deploy: sync /opt/clonr to origin/main and rebuild
# WARNING: This unit stomps local uncommitted changes on this host.
# If you need to test uncommitted work, stop the timer first:
#   systemctl stop clonr-autodeploy.timer
# Resume with:
#   systemctl start clonr-autodeploy.timer
Documentation=https://github.com/sqoia-dev/clonr/blob/main/scripts/autodeploy/
After=network-online.target
Wants=network-online.target
# Prevent concurrent runs if a previous build is still running
After=clonr-serverd.service

[Service]
Type=oneshot
# Script is installed to /usr/local/sbin/ (NOT the repo path) so it survives
# git reset --hard operations. The install.sh copies it there. When the script
# updates itself (via git pull), the new version takes effect on the next cycle
# because install.sh is re-run as part of any manual upgrade path.
#
# NOTE: The script in /opt/clonr/ is the source of truth, but the installed
# copy at /usr/local/sbin/clonr-autodeploy.sh is what runs. If the script
# logic changes, re-run scripts/autodeploy/install.sh to update the installed copy.
ExecStart=/usr/local/sbin/clonr-autodeploy.sh
# Run as root: needs to write to /usr/local/bin, restart clonr-serverd,
# and write to /var/lib/clonr (owned by clonr user but root can write)
User=root

# Output goes to journal — tail with:
#   journalctl -u clonr-autodeploy.service -f
StandardOutput=journal
StandardError=journal

# Allow a generous timeout for builds on this lab box (initramfs build
# includes cpio/gzip which can take 30-60s on slow storage)
TimeoutStartSec=10min

# Do not restart on failure — the timer will fire again in 2 minutes.
# Exit non-zero is used to signal "retry next cycle" intentionally.
Restart=no
