# deploy/systemd/clustr-autodeploy.timer
[Unit]
Description=clustr auto-deploy timer — poll origin/main every 2 minutes
Documentation=https://github.com/sqoia-dev/clustr/blob/main/scripts/autodeploy/

[Timer]
# Run 1 minute after boot (allow network to be fully up)
OnBootSec=1min
# Then every 2 minutes, measured from when the service last *activated*
# (i.e., started, not finished) — gives ~2 min cycles even if build takes 90s
OnUnitActiveSec=2min
# Persist=true: if the host was off when a timer would have fired,
# catch up once on next boot (handles overnight reboots gracefully)
Persistent=true
# Reduce stampede when multiple units fire at the same time
RandomizedDelaySec=10s

[Install]
WantedBy=timers.target
