# deploy/ansible/roles/clustr/templates/clustr-serverd.service.j2
# clustr-serverd systemd unit — generated by Ansible
# Managed by Ansible role clustr — do not edit manually.
# Re-run the playbook to update this file.

[Unit]
Description=clustr Provisioning Server
After=network.target dev-pts.mount
Wants=network.target dev-pts.mount

[Service]
Type=simple

ExecStartPre=/bin/bash -c "mountpoint -q /dev/pts || mount -t devpts devpts /dev/pts -o ptmxmode=666,mode=620,gid=5"
ExecStart={{ clustr_bin_dir }}/clustr-serverd --pxe

# ---------------------------------------------------------------------------
# Environment files
# ---------------------------------------------------------------------------
# clustr.env — non-secret runtime config
EnvironmentFile={{ clustr_config_dir }}/clustr.env
# secrets.env — CLUSTR_SECRET_KEY, CLUSTR_SESSION_SECRET
# The '-' prefix means systemd will not fail if the file does not exist yet.
EnvironmentFile=-{{ clustr_config_dir }}/secrets.env

# ---------------------------------------------------------------------------
# Restart policy
# ---------------------------------------------------------------------------
Restart=on-failure
RestartSec=5s
# Allow up to 300s for graceful shutdown — ISO builds take 10-15 minutes.
TimeoutStopSec=300

# ---------------------------------------------------------------------------
# Resource limits
# ---------------------------------------------------------------------------
TasksMax=infinity
MemoryMax=8G
OOMScoreAdjust=-500

# ---------------------------------------------------------------------------
# Security hardening
# ---------------------------------------------------------------------------
# See deploy/systemd/clustr-serverd.service for detailed rationale comments
# explaining why each directive is set the way it is.
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false
RestrictNamespaces=false
MountFlags=

# ---------------------------------------------------------------------------
# Capabilities
# ---------------------------------------------------------------------------
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN CAP_SYS_CHROOT CAP_MKNOD CAP_SETPCAP CAP_AUDIT_WRITE CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN CAP_SYS_CHROOT CAP_MKNOD CAP_SETPCAP CAP_AUDIT_WRITE CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_SETUID CAP_SETGID

# ---------------------------------------------------------------------------
# Device access
# ---------------------------------------------------------------------------
DevicePolicy=auto
DeviceAllow=/dev/kvm rwm
DeviceAllow=/dev/net/tun rwm
DeviceAllow=block-loop rwm
DeviceAllow=/dev/loop-control rw

[Install]
WantedBy=multi-user.target
